# docker-dev/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:mysql]
command=/usr/bin/mysqld_safe --defaults-file=/etc/my.cnf --user=mysql --datadir=/var/lib/mysql --bind-address=0.0.0.0 --socket=/var/run/mysqld/mysqld.sock
autostart=true
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql.log
user=mysql
priority=100
environment=HOME="/var/lib/mysql"

[program:redis]
command=/usr/bin/redis-server /etc/redis.conf --daemonize no
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log
user=redis
priority=200
directory=/var/lib/redis

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm.log
user=root
priority=300
# entrypoint 脚本会检测 PHP-FPM 路径并更新此命令

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
user=root
priority=400

[program:nuxt]
command=/bin/sh -c "sleep 5; cd /data/httpd/ECShopX_desktop-frontend && if [ -d .nuxt ] || [ -d .output ]; then if npm run 2>/dev/null | grep -q 'start'; then PORT=3000 HOST=0.0.0.0 npm run start; else PORT=3000 HOST=0.0.0.0 npm run dev; fi; else exit 0; fi"
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/nuxt.log
stderr_logfile=/var/log/supervisor/nuxt.log
user=root
priority=500
directory=/data/httpd/ECShopX_desktop-frontend

[include]
files = /etc/supervisor/conf.d/*.conf
